# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:proftpd]
command=/usr/sbin/proftpd 

